# langin_page.conf
[program:landing-page-5000]
command = python /home/ubuntu/campus_vault_landing/manage.py
user = ubuntu
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8